# board/garimpo/rootfs_overlay/etc/network/if-pre-up.d/wpasupplicant
#!/bin/sh

if [ "$IFACE" = wlan0 ]; then
    wpa_supplicant -B -c/etc/wpa_supplicant.conf -iwlan0
fi
